# twitch_bot.service
[Unit]
Description=Twitch_bot

[Service]
ExecStart=/usr/bin/python3 /app/bot.py
Restart=always

[Install]
WantedBy=multi-user.target
